# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv sh

exec \
  s6-setuidgid transmission \
    /usr/bin/transmission-daemon -g /config -c /watch -f
